# workers.conf
[program:register]
command=/usr/bin/php -f registerWorker.php
process_name=%(program_name)s_%(process_num)02d
numprocs=25
autorestart=true
user=vagrant
directory=/home/vagrant/code/app
startretries=5